# deploy/openclaw-backend.service
[Unit]
Description=OpenClaw Backend (FastAPI/uvicorn)
After=network.target
Wants=network.target

[Service]
Type=simple
User=zenca
WorkingDirectory=/home/zenca/openclaw-core
EnvironmentFile=-/home/zenca/openclaw-core/.env
ExecStart=/home/zenca/openclaw-core/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-backend

[Install]
WantedBy=multi-user.target
